# supervisor.conf
[unix_http_server]
file=/tmp/supervisord/supervisor.sock ; (the path to the socket file)
chmod=0700                  ; socket file mode (default 0700)
chown=umail:umail           ; socket file uid:gid owner

[supervisord]
logfile=/tmp/supervisord/supervisord.log ; (main log file;default $CWD/supervisord.log)
logfile_maxbytes=50MB       ; (max main logfile bytes b4 rotation;default 50MB)
logfile_backups=5           ; (num of main logfile rotation backups;default 10)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)
user=root                  ; (default is current user, required if root)
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisord/supervisor.sock; use a unix:// URL  for a unix socket
prompt=root               ; cmd line prompt (default "supervisor")

[program:app]
command                 = /home/python/pyenv/versions/flask-blog/bin/gunicorn -c gconf.py manage:app
user                    = root
directory               = /home/python/git/flask-blog
autostart               = true
autorestart             = true
redirect_stderr         = true
stdout_logfile          = /home/python/git/flask-blog/log/error.log
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups  = 10